# deploy/com.personalwiki.site.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!-- Copy to ~/Library/LaunchAgents/com.personalwiki.site.plist (drop .example),
     replace /ABS/PATH/TO/personal_wiki with this repo's absolute path, then:
     launchctl load ~/Library/LaunchAgents/com.personalwiki.site.plist -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.personalwiki.site</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/ABS/PATH/TO/personal_wiki/scripts/serve.mjs</string>
    <string>--host</string><string>127.0.0.1</string><string>--port</string><string>4321</string>
  </array>
  <key>WorkingDirectory</key><string>/ABS/PATH/TO/personal_wiki</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/tmp/personalwiki-site.log</string>
  <key>StandardErrorPath</key><string>/tmp/personalwiki-site.err</string>
</dict>
</plist>
